# deployment/s6-overlay/s6-rc.d/bind/run
#!/bin/bash
test -f /etc/default/named && . /etc/default/named
exec /usr/sbin/named -u bind -f