# deployment/s6-overlay/s6-rc.d/bind/run
#!/command/execlineb -P
/usr/sbin/named -u named -g