# deployment/s6-overlay/s6-rc.d/dyndns/run
#!/command/execlineb -P
with-contenv
cd /root
/root/dyndns
